# vendor/bundle/ruby/3.4.0/gems/foreman-0.87.2/data/export/systemd/master.target.erb
[Unit]
Wants=<%= service_names.join(' ') %>

[Install]
WantedBy=multi-user.target
