# deployment/framky/chatwoot.target
[Unit]
Description=Chatwoot Application (crm.framky.com)
Wants=chatwoot-web.service chatwoot-worker.service

[Install]
WantedBy=multi-user.target
