# vendor/bundle/ruby/3.4.0/gems/foreman-0.87.2/spec/resources/export/supervisord/app-alpha-2.conf
[program:app-alpha-1]
command=./alpha
autostart=true
autorestart=true
stdout_logfile=/var/log/app/alpha-1.log
stderr_logfile=/var/log/app/alpha-1.error.log
user=app
directory=/tmp/app
environment=PORT="5000"

[program:app-alpha-2]
command=./alpha
autostart=true
autorestart=true
stdout_logfile=/var/log/app/alpha-2.log
stderr_logfile=/var/log/app/alpha-2.error.log
user=app
directory=/tmp/app
environment=PORT="5001"

[group:app]
programs=app-alpha-1,app-alpha-2
