# .conflict-side-0/deployment/chatwoot.target
[Unit]
Wants=chatwoot-web.1.service chatwoot-worker.1.service

[Install]
WantedBy=multi-user.target