# vendor/bundle/ruby/3.4.0/gems/foreman-0.87.2/spec/resources/export/upstart/app-bravo-1.conf
start on starting app-bravo
stop on stopping app-bravo
respawn

env PORT=5100

setuid app

chdir /tmp/app

exec ./bravo
