# vendor/bundle/ruby/3.4.0/gems/foreman-0.87.2/spec/resources/export/upstart/app-alpha-2.conf
start on starting app-alpha
stop on stopping app-alpha
respawn

env PORT=5001

setuid app

chdir /tmp/app

exec ./alpha
