# vendor/bundle/ruby/3.4.0/gems/foreman-0.87.2/spec/resources/export/upstart/app-bravo.conf
start on starting app
stop on stopping app
