# vendor/bundle/ruby/3.4.0/gems/foreman-0.87.2/spec/resources/export/upstart/app.conf
start on runlevel [2345]
stop on runlevel [!2345]
