# vendor/bundle/ruby/3.4.0/gems/foreman-0.87.2/spec/resources/export/daemon/app-bravo-1.conf
start on starting app-bravo
stop on stopping app-bravo
respawn

env PORT=5100

exec start-stop-daemon --start --chuid app --chdir /tmp/app --make-pidfile --pidfile /var/run/app/app-bravo-1.pid --exec ./bravo >> /var/log/app/app-bravo-1.log 2>&1
