# vendor/bundle/ruby/3.4.0/gems/foreman-0.87.2/spec/resources/export/daemon/app-alpha-2.conf
start on starting app-alpha
stop on stopping app-alpha
respawn

env PORT=5001

exec start-stop-daemon --start --chuid app --chdir /tmp/app --make-pidfile --pidfile /var/run/app/app-alpha-2.pid --exec ./alpha >> /var/log/app/app-alpha-2.log 2>&1
