# vendor/bundle/ruby/3.4.0/gems/foreman-0.87.2/spec/resources/export/daemon/app.conf
pre-start script

bash << "EOF"
mkdir -p /var/log/app
chown -R app /var/log/app
mkdir -p /var/run/app
chown -R app /var/run/app
EOF

end script

start on runlevel [2345]

stop on runlevel [016]
