# vendor/bundle/ruby/3.4.0/gems/foreman-0.87.2/spec/resources/export/systemd/app-alpha.1.service
[Unit]
PartOf=app.target
StopWhenUnneeded=yes

[Service]
User=app
WorkingDirectory=/tmp/app
Environment=PORT=5000
Environment=PS=alpha.1
ExecStart=/bin/bash -lc 'exec -a "app-alpha.1" ./alpha'
Restart=always
RestartSec=14s
StandardInput=null
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=%n
KillMode=mixed
TimeoutStopSec=5
