# vendor/bundle/ruby/3.4.0/gems/foreman-0.87.2/spec/resources/export/systemd/app.target
[Unit]
Wants=app-alpha.1.service app-bravo.1.service app-foo_bar.1.service app-foo-bar.1.service

[Install]
WantedBy=multi-user.target
