# vendor/bundle/ruby/3.4.0/gems/foreman-0.87.2/spec/resources/export/runit/app-bravo-1/run
#!/bin/sh
cd /tmp/app
exec 2>&1
exec chpst -u app -e /tmp/init/app-bravo-1/env ./bravo
